# deploy/hl-detect-signals.timer
[Unit]
Description=HL Signals — detect consensus signals (every 3 min)

[Timer]
OnBootSec=3min
OnUnitActiveSec=3min

[Install]
WantedBy=timers.target
