# deploy/hl-refresh-positions.service
[Unit]
Description=HL Signals — refresh open positions
After=network-online.target

[Service]
Type=oneshot
User=www-data
WorkingDirectory=/var/www/hyperliquid-signals
ExecStart=/usr/bin/flock -n /var/www/hyperliquid-signals/data/hl-refresh-positions.lock /usr/bin/php /var/www/hyperliquid-signals/cron/refresh_positions.php
Nice=5
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
ReadWritePaths=/var/www/hyperliquid-signals/data
