# deploy/hl-refresh-positions.timer
[Unit]
Description=HL Signals — refresh positions (every 30s)

[Timer]
OnBootSec=30s
OnUnitActiveSec=30s
AccuracySec=5s

[Install]
WantedBy=timers.target
