# deploy/hl-refresh-positions.service
[Unit]
Description=HL Signals — refresh open positions
After=network-online.target

[Service]
Type=oneshot
User=www-data
WorkingDirectory=/var/www/hyperliquid-signals
ExecStart=/usr/bin/php /var/www/hyperliquid-signals/cron/refresh_positions.php
Nice=5
