# deploy/hl-refresh-portfolios.service
[Unit]
Description=HL Signals — refresh portfolios
After=network-online.target

[Service]
Type=oneshot
User=www-data
WorkingDirectory=/var/www/hyperliquid-signals
ExecStart=/usr/bin/flock -n /run/hl-refresh-portfolios.lock /usr/bin/php /var/www/hyperliquid-signals/cron/refresh_portfolios.php
Nice=10
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
ReadWritePaths=/var/www/hyperliquid-signals/data /run
RuntimeDirectory=hl-refresh-portfolios
