# deploy/hl-seed-traders.service
[Unit]
Description=HL Signals — discover new traders via recentTrades
After=network-online.target

[Service]
Type=oneshot
User=www-data
WorkingDirectory=/var/www/hyperliquid-signals
ExecStart=/usr/bin/flock -n /run/hl-seed-traders.lock /usr/bin/php /var/www/hyperliquid-signals/cron/seed_traders.php
Nice=15
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
ReadWritePaths=/var/www/hyperliquid-signals/data /run
